# services/systemd-user/professional/chek-teleop-adapter.service
[Unit]
Description=CHEK Isaac Teleop Adapter
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=__EDGE_ROOT__
Environment=CHEK_EDGE_ROOT=__CHEK_EDGE_ROOT__
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env bash __EDGE_ROOT__/scripts/start_isaac_host_adapter.sh
Restart=always
RestartSec=2
RestartPreventExitStatus=78
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=default.target
