# services/systemd/professional/chek-edge-wifi-sensing.service
[Unit]
Description=CHEK Edge Wi-Fi DensePose Sensing Server
After=chek-edge-stack.service
Wants=chek-edge-stack.service

[Service]
Type=simple
WorkingDirectory=__EDGE_ROOT__
Environment=CHEK_EDGE_ROOT=__CHEK_EDGE_ROOT__
ExecStart=/usr/bin/env bash __EDGE_ROOT__/scripts/edge_wifi_sensing_autostart.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
